# conf/wake.cron
*/1 * * * 6,7 echo "dan" >> /workspace/smart-wake/conf/wake.dan #ma beat
