# conf/wake.cron
00 7 * * 1,2,3,4,5 /root/proq/scripts/wakeup.py #Cours
30 9 * * 6,7 /root/proq/scripts/wakeup.py #Week'end
#30 10 * * 1,2,3,4,5,6,7 /root/proq/scripts/wakeup.py #Vacances
